# vita.service
[Service]
ExecStart=/bin/bash /home/pi/scripts/launch.sh